# ops/systemd/agendasalon-email.service
[Unit]
Description=AgendaSalon transactional email queue
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=agendasalon
Group=www-data
WorkingDirectory=/var/www/agendasalon/app
EnvironmentFile=/etc/agendasalon/agendasalon.env
ExecStart=/var/www/agendasalon/.venv/bin/python /var/www/agendasalon/app/manage.py process_outbound_emails --limit 100
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/agendasalon/shared
