# ops/systemd/agendasalon-demo-refresh-dispatch.timer
[Unit]
Description=Check for a requested AgendaSalon demo refresh every 30 seconds

[Timer]
OnBootSec=30s
OnUnitInactiveSec=30s
Persistent=false
AccuracySec=1s
RandomizedDelaySec=0
Unit=agendasalon-demo-refresh-dispatch.service

[Install]
WantedBy=timers.target
