# ops/systemd/backup-agendasalon.service
[Unit]
Description=AgendaSalon local authenticated backup
After=postgresql.service
Requires=postgresql.service

[Service]
Type=oneshot
User=agendasalon
Group=www-data
WorkingDirectory=/var/www/agendasalon/app
EnvironmentFile=/etc/agendasalon/agendasalon.env
UMask=0077
ExecStart=/var/www/agendasalon/.venv/bin/python /var/www/agendasalon/app/manage.py backup_agendasalon --backup-root /var/backups/agendasalon --media-root /var/www/agendasalon/shared/media --destination local
ExecStartPost=/var/www/agendasalon/.venv/bin/python -m ops.backup_restore retention --backup-root /var/backups/agendasalon --daily 7 --weekly 4 --monthly 6 --apply
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
ReadOnlyPaths=/var/www/agendasalon/app /var/www/agendasalon/shared/media
ReadWritePaths=/var/backups/agendasalon
