# ops/systemd/backup-agendasalon.timer
[Unit]
Description=Daily AgendaSalon local backup

[Timer]
OnCalendar=*-*-* 03:20:00
Persistent=true
RandomizedDelaySec=15m
Unit=backup-agendasalon.service

[Install]
WantedBy=timers.target
