# ops/systemd/agendasalon-registration-purge.service.d/10-demo-refresh-safety.conf
[Service]
ExecStartPre=/usr/local/libexec/agendasalon-demo-start-guard
